# scripts/systemd/0xkobold.service
[Unit]
Description=0xKobold Multi-Agent Automation Platform
After=network.target
Wants=network.target

[Service]
Type=simple
User=%I
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="HOME=/home/%I"
Environment="NODE_ENV=production"

# Change this to match your bun installation
ExecStart=/usr/local/bin/bun run /usr/local/bin/0xkobold daemon
ExecStop=/usr/local/bin/bun /usr/local/bin/0xkobold stop
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.0xkobold %h/.config

[Install]
WantedBy=multi-user.target
